# crontab
0 2 * * * root /renew.sh > /dev/stdout 2>&1
