# supervisord.conf
[supervisord]
nodaemon  = true
environment = PLACEHOLDER=true

[unix_http_server]
file = /supervisor.sock
chmod = 0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///supervisor.sock

[program:cron]
command   = /usr/bin/go-cron -file="/etc/crontab"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command   = /usr/sbin/nginx
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
